# app/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/apache2 -DFOREGROUND
autostart=true
autorestart=true

